# scripts/ops/showtracker-schedule-confidence.timer
[Unit]
Description=Run ShowTracker schedule confidence reconciler nightly

[Timer]
OnCalendar=*-*-* 03:20:00 UTC
RandomizedDelaySec=10m
Persistent=true
Unit=showtracker-schedule-confidence.service

[Install]
WantedBy=timers.target
